# src/scripts/bus_env_cbus/boot_scripts/systemd/system/rosmasterd.service
[Unit]
Description=ROS master service
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/rosmasterd.sh

[Install]
WantedBy=default.target
